# deploy/husky.service
[Unit]
Description=Husky — AI assistant service
After=network.target

[Service]
Type=simple
User=__HUSKY_USER__
WorkingDirectory=__HUSKY_INSTALL_DIR__
EnvironmentFile=__HUSKY_ENV_FILE__
ExecStart=/usr/bin/java -jar __HUSKY_SERVICE_JAR__
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=__HUSKY_DATA_DIR__ __HUSKY_DATA_DIR__/config __HUSKY_DATA_DIR__/skills __HUSKY_DATA_DIR__/db __HUSKY_DATA_DIR__/logs

[Install]
WantedBy=multi-user.target
